# scripts/forged.service
[Unit]
Description=Forge node daemon
After=network.target

[Service]
Type=simple
User=forge
Group=forge
ExecStart=/usr/local/bin/forged
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
